# docker/sv-crond.conf
[program:cron]
autorestart=false
command=/usr/sbin/cron -f
